# services/chromium-kiosk.service
[Unit]
Description=Chromium Kiosk Mode Display (HDMI)
After=graphical.target network-online.target
Wants=graphical.target network-online.target

[Service]
Type=simple
User=pi
Group=pi
Environment=DISPLAY=:0
Environment=WAYLAND_DISPLAY=wayland-0
Environment=XDG_RUNTIME_DIR=/run/user/1000
Environment=XAUTHORITY=/home/pi/.Xauthority
ExecStartPre=/home/pi/pull_tournament_data.sh
ExecStart=/home/pi/chromium_kiosk.sh
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=graphical.target
