# services/tournament-data-boot.service
[Unit]
Description=Pull Tournament Data at Boot
After=network-online.target
Wants=network-online.target
Before=chromium-kiosk.service

[Service]
Type=oneshot
User=pi
Group=pi
ExecStart=/home/pi/pull_tournament_data.sh
RemainAfterExit=yes
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
